# linux-installer/upstart/li_metrics_agent.conf
description "Load Impact Agent for server metrics"
author      "Load Impact AB"

start on filesystem or runlevel [2345]
stop on shutdown

respawn

exec /usr/bin/li-metrics-agent-run
